# roles/user/templates/user.service
[Unit]
Description = User Service
[Service]
User=roboshop
Environment=MONGO=true
Environment=REDIS_HOST=redis-{{env}}.rdevopsb73.online
Environment=MONGO_URL="mongodb://mongodb-{{env}}.rdevopsb73.online:27017/users"
ExecStart=/bin/node /app/server.js
SyslogIdentifier=user

[Install]
WantedBy=multi-user.target
